# deploy/sketchymine-deployment.path
[Unit]
Description=SketchyMine deployment path monitor

[Path]
PathModified=/home/gitlab/sketchymine/deployment.raw
Unit=sketchymine-deployment.service

[Install]
WantedBy=multi-user.target