# deploy/sketchymine-deployment.service
[Unit]
Description=SketchyMine deployment service

[Service]
ExecStart=/home/gitlab/sketchymine/deploy.sh